# 99-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/home/praani/.config/qtile/monitor-hotplug.sh"
